# scripts/crontab
SHELL=/bin/bash
PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin

# Run database backup daily at 2:00 AM
0 2 * * * root /backups/backup-db.sh
